# services/backend.service
[Unit]
Description=Fundvís Backend Service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
User=raspberrypi
Group=raspberrypi
WorkingDirectory=/opt/properties-by-magni/backend
EnvironmentFile=/opt/properties-by-magni/.env
Environment="PATH=/opt/properties-by-magni/backend/venv/bin"
ExecStart=/opt/properties-by-magni/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=append:/var/log/backend.log
StandardError=append:/var/log/backend.log
SyslogIdentifier=properties-backend

[Install]
WantedBy=multi-user.target
