# services/scraper.service
[Unit]
Description=Property Scraper Service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
User=raspberrypi
Group=raspberrypi
WorkingDirectory=/opt/properties-by-magni/scraper
Environment="PATH=/opt/properties-by-magni/scraper/venv/bin"
ExecStart=/opt/properties-by-magni/scraper/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/scraper.log
StandardError=append:/var/log/scraper.log
SyslogIdentifier=scraper

[Install]
WantedBy=multi-user.target
